# battery_sync/crontab
*/5 *   * * *   root    cd /data/data/com.termux/files/home/battery_sync/ && ./battery_sync.sh
